# crontab
0 1 * * * curl http://workhours-scrapper.flycast/api/scrape