# crontab
0 3 * * * curl http://workhours-scrapper.flycast/api/scrape